# systemd/stave-synth.service.d/faust.conf
# Faust-native DSP paths. Shipped 2026-04-18 after user-approved A/B on all modules.
# To disable any single module, override with:
#   systemctl --user edit stave-synth
# and set e.g. `Environment=STAVE_FAUST_REVERB=0`.
[Service]
Environment=STAVE_FAUST_REVERB=1
Environment=STAVE_FAUST_PING_PONG=1
Environment=STAVE_FAUST_OSC_BANK=1
Environment=STAVE_FAUST_SYMPATHETIC=1
Environment=STAVE_FAUST_MASTER_FX=1
Environment=STAVE_FAUST_BUS_COMP=1
Environment=STAVE_FAUST_ORGAN=1
